# workspaces/worker/crontab
* * * * * bun /worker/src/everyMinute.ts 2>&1
0 * * * * bun /worker/src/everyHour.ts 2>&1
0 0 * * * bun /worker/src/everyDay.ts 2>&1